# xdg-autostart/.config/autostart/bar.desktop
[Desktop Entry]
Type=Application
Name=bar.sh
Exec=/home/jordif/bin/bar.sh
OnlyShowIn=i3
