# xdg-autostart/.config/autostart/battray.desktop
[Desktop Entry]
Name=Battery Tray
Icon=/usr/share/battray/icon.png
Exec=battray
Terminal=false
Type=Application
NoDisplay=true
OnlyShowIn=qtile;
