# xdg-autostart/.config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=dunst
Exec=dunst
OnlyShowIn=i3;bspwm;
