# xdg-autostart/.config/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=firefox
Exec=firefox
