# xdg-autostart/.config/autostart/google-chrome.desktop
[Desktop Entry]
Type=Application
Name=Google Chrome
Exec=google-chrome
