# xdg-autostart/.config/autostart/keepassx.desktop
[Desktop Entry]
Type=Application
Name=keepassx2
Exec=keepassx2
