# xdg-autostart/.config/autostart/keepassxc.desktop
[Desktop Entry]
Type=Application
Name=keepassxc
Exec=keepassxc
