# xdg-autostart/.config/autostart/nitrogen.desktop
[Desktop Entry]
Type=Application
Name=nitrogen
Exec=nitrogen --restore
OnlyShowIn=i3;bspwm;
