# xdg-autostart/.config/autostart/obsidian.desktop
[Desktop Entry]
Type=Application
Name=Slack
Exec=flatpak run md.obsidian.Obsidian
