# xdg-autostart/.config/autostart/playerctld.desktop
[Desktop Entry]
Type=Application
Name=playerctld
Exec=playerctld daemon
