# xdg-autostart/.config/autostart/powermate.desktop
[Desktop Entry]
Type=Application
Name=sxhkd
Exec=$HOME/code/powermater/powermate.sh
