# xdg-autostart/.config/autostart/powermate.desktop
[Desktop Entry]
Type=Application
Name=powermate
Exec=/home/jordif/code/powermater/powermate.sh
