# xdg-autostart/.config/autostart/slack.desktop
[Desktop Entry]
Type=Application
Name=Slack
Exec=flatpak run com.slack.Slack
