# xdg-autostart/.config/autostart/slack.desktop
[Desktop Entry]
Name=Slack
StartupWMClass=Slack
Exec=flatpak run com.slack.Slack
Icon=slack
Type=Application
Categories=Network;InstantMessaging;
