# xdg-autostart/.config/autostart/spotify.desktop
[Desktop Entry]
Type=Application
Name=spotify
Exec=flatpak run com.spotify.Client
