# xdg-autostart/.config/autostart/sxhkd.qtile.desktop
[Desktop Entry]
Type=Application
OnlyShowIn=qtile
Name=sxhkd
Exec=/usr/bin/sxhkd -c /home/jordif/.config/sxhkd/sxhkdrc.common
