# xdg-autostart/.config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Name=sxhkd
Exec=/usr/bin/sxhkd
