# xdg-autostart/.config/autostart/telegram.desktop
[Desktop Entry]
Type=Application
Name=telegram
Exec=flatpak run org.telegram.desktop
NotShowIn=qtile
