# xdg-autostart/.config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Name=udiskie
Exec=udiskie --tray
OnlyShowIn=i3;qtile;
