# kmonad/.config/autostart/kmonad.desktop
[Desktop Entry]
Type=Application
Name=kmonad
Exec=kmonad /home/jordif/.config/kmonad/custom_iso_tkl.kbd
Hidden=true
