# kmonad/.config/systemd/user/kmonad.service
[Unit]
Description=kmonad

[Service]
ExecStart=/home/jordif/.local/bin/kmonad /home/jordif/.config/kmonad/laptop.kbd

[Install]
WantedBy=default.target
