# systemd/.config/systemd/user/fusuma.service
[Unit]
Description=Fusuma

[Service]
ExecStart=/home/jordif/bin/fusuma

[Install]
WantedBy=default.target
